# ops/agency-crm.service
[Unit]
Description=Agency CRM (Node)
After=network.target

[Service]
Type=simple
User=xipx
WorkingDirectory=/opt/agency-crm
EnvironmentFile=/opt/agency-crm/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/env node /opt/agency-crm/server/server.js
Restart=on-failure
RestartSec=3

# Hardening (safe defaults)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target
